# files/clamdscan.timer
[Unit]
Description=Scan whole system using ClamAV

[Timer]
OnCalendar=23:00:00
Persistent=true

[Install]
WantedBy=timers.target